# deployment/junior_assistant.service
[Unit]
Description=Júnior Affiliate Assistant API Service
After=network.target

[Service]
User=dev
Group=www-data
WorkingDirectory=/home/dev/junior-affiliate-assistant
ExecStart=/home/dev/junior-affiliate-assistant/venv/bin/uvicorn main:app --host 0.0.0.0 --port 8000
Restart=always

[Install]
WantedBy=multi-user.target